# confs/results.conf
# description "uWSGI server for {{ PROJECT_SLUG }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    . /etc/environment
    /usr/bin/python {{ SERVER_REPOSITORY_PATH }}/daemon/daemon.py
end script
