# confs/results.conf
# description "uWSGI server for {{ PROJECT_SLUG }}"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    . /etc/environment
    {{ SERVER_REPOSITORY_PATH }}/{{ PROJECT_SLUG }}/
end script
